# daily_call_data/crontab
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

1 4 * * /call_data/run_quote.sh >> /call_data/cron.log 2>&1
